# compose/dev/start_nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
user=nginx
stderr_logfile=/var/log/nginx_error.log
stdout_logfile=/var/log/nginx_out.log